# deploy/levelprobs.service
[Unit]
Description=levelprobs live SPX level-hit alerts (Yahoo Finance ^GSPC)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/home/matt/git/levelprobs
# Real, live data via Yahoo Finance (^GSPC 5-min bars). No account, no tunnel:
# a plain outbound poll. History is pulled from Yahoo automatically, and alerts
# are suppressed while the market is closed.
ExecStart=/usr/bin/python3 -m levelprobs SPX --watch --feed yahoo \
    --poll 20 --min-prob 0.80 --min-points 15
Restart=always
RestartSec=5
# desktop notifications need the user session bus (inherited by --user services)

[Install]
WantedBy=graphical-session.target
